# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/prometheus/run
#!/usr/bin/with-contenv bashio

echo "Starting Prometheus..."
sleep 3
exec /opt/prometheus/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/data/prometheus \
  --web.console.libraries=/opt/prometheus/console_libraries \
  --web.console.templates=/opt/prometheus/consoles \
  --web.enable-lifecycle 