# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/prometheus/run
#!/command/with-contenv bashio

# Start Prometheus
exec s6-setuidgid root /opt/prometheus/prometheus \
    --config.file=/etc/prometheus/prometheus.yml \
    --storage.tsdb.path=/data/prometheus \
    --web.listen-address=:9090 \
    --web.enable-lifecycle \
    --web.enable-admin-api 