# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/prometheus/up
#!/usr/bin/with-contenv bashio

# Wait for Prometheus to be ready
for _ in $(seq 1 30); do
    if curl -s "http://localhost:9090/-/ready" > /dev/null; then
        echo "Prometheus is ready"
        exit 0
    fi
    sleep 0.5
done

echo "Prometheus failed to become ready"
exit 1 