# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/prometheus/up
#!/usr/bin/with-contenv bashio

echo "Prometheus service is up" 