# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/prometheus/up
#!/command/with-contenv bashio

# Wait for Prometheus to be ready
for i in $(seq 1 30); do
    if curl -s "http://localhost:9090/-/ready" > /dev/null; then
        echo "Prometheus is ready"
        exit 0
    fi
    echo "Waiting for Prometheus... (attempt $i/30)"
    sleep 0.5
done

echo "Prometheus failed to become ready after 15 seconds"
exit 1 