# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
echo "Starting NGINX for ingress routing..."
exec nginx 