# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/blackbox-exporter/run
#!/command/with-contenv bashio

# Start blackbox exporter
exec s6-setuidgid root /opt/blackbox_exporter/blackbox_exporter \
    --config.file=/etc/blackbox_exporter/blackbox.yml \
    --web.listen-address=:9115 