# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/blackbox-exporter/up
#!/usr/bin/with-contenv bashio

echo "Blackbox Exporter service is up" 