# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/blackbox-exporter/up
#!/command/with-contenv bashio

# Wait for Blackbox Exporter to be ready
for _ in $(seq 1 30); do
    if curl -s "http://localhost:9115/health" > /dev/null; then
        echo "Blackbox Exporter is ready"
        exit 0
    fi
    sleep 0.5
done

echo "Blackbox Exporter failed to become ready"
exit 1 