# influxdb-stack/rootfs/etc/s6-overlay/s6-rc.d/grafana/up
#!/command/with-contenv bashio

# Wait for Grafana to be ready
for i in {1..120}; do
    if curl -f http://localhost:3000/api/health > /dev/null 2>&1; then
        echo "Grafana is ready"
        exit 0
    fi
    echo "Waiting for Grafana to start..."
    sleep 0.5
done

echo "Grafana failed to start within timeout"
exit 1 