# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/vscode-api/run
#!/command/with-contenv bashio

echo "Starting VS Code API server..."
exec s6-setuidgid root /usr/local/bin/vscode-api-server 