# influxdb-stack/rootfs/etc/s6-overlay/s6-rc.d/influxdb/run
#!/command/with-contenv bashio

echo "Starting InfluxDB..."
exec s6-setuidgid root /opt/influxdb/usr/bin/influxd \
    --store=disk \
    --engine-path=/data/influxdb \
    --http-bind-address=:8086 \
    --reporting-disabled 