# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/alertmanager/run
#!/command/with-contenv bashio

echo "Starting Alertmanager..."
exec s6-setuidgid root /opt/alertmanager/alertmanager \
    --config.file=/etc/alertmanager/alertmanager.yml \
    --storage.path=/data/alertmanager \
    --web.listen-address=:9093 