# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/alertmanager/run
#!/usr/bin/with-contenv bashio

echo "Starting Alertmanager..."
sleep 2
exec /opt/alertmanager/alertmanager \
  --config.file=/etc/alertmanager/alertmanager.yml \
  --storage.path=/data/alertmanager 