# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/alertmanager/up
#!/command/with-contenv bashio

# Wait for Alertmanager to be ready
for i in $(seq 1 30); do
    if curl -s "http://localhost:9093/-/ready" > /dev/null; then
        echo "Alertmanager health endpoint is ready"
        
        # Also verify API endpoint is ready (what Karma actually needs)
        if curl -s "http://localhost:9093/api/v2/status" > /dev/null 2>&1; then
            echo "Alertmanager API is ready"
            exit 0
        else
            echo "Alertmanager health ready but API not ready yet (attempt $i/30)"
        fi
    fi
    echo "Waiting for Alertmanager... (attempt $i/30)"
    sleep 0.5
done

echo "Alertmanager failed to become ready after 15 seconds"
exit 1 