# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/alertmanager/up
#!/usr/bin/with-contenv bashio

echo "Alertmanager service is up" 