# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/promtail/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Promtail
# ==============================================================================

bashio::log.info "Starting Promtail..."

exec s6-setuidgid promtail promtail -config.file=/etc/promtail/promtail.yml 