# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/code-server/run
#!/command/with-contenv bashio

echo "Starting VS Code Server..."

# Always use /etc as the workspace
WORKSPACE="/etc"

# Ensure .vscode/settings.json is present in the workspace
if [ ! -d "$WORKSPACE/.vscode" ]; then
  mkdir -p "$WORKSPACE/.vscode"
fi
if [ -f /etc/.vscode/settings.json ]; then
  cp -f /etc/.vscode/settings.json "$WORKSPACE/.vscode/settings.json"
fi

# Create workspace directory if it doesn't exist
mkdir -p "$WORKSPACE"

# Fix Node.js binary issue by patching the code-server script
echo "Fixing code-server to use system Node.js..."
sed -i 's|exec "$ROOT/lib/node" "$ROOT" "$@"|exec "/usr/bin/node" "$ROOT" "$@"|g' /opt/code-server/bin/code-server

# Set environment variables for code-server
export PORT=8443
export HOST=0.0.0.0
export USER_DATA_DIR=/opt/code-server/data
export EXTENSIONS_DIR=/opt/code-server/extensions

# Start code-server without authentication
exec s6-setuidgid root /opt/code-server/bin/code-server \
    --bind-addr 0.0.0.0:8443 \
    --user-data-dir /opt/code-server/data \
    --extensions-dir /opt/code-server/extensions \
    --auth none \
    --disable-telemetry \
    --disable-update-check \
    "$WORKSPACE" 