# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/loki/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Loki..."

# Wait for filesystem to be ready
sleep 2

exec /usr/local/bin/loki \
    -config.file=/etc/loki/loki.yml \
    -log.level=info \
    >> /var/log/loki/loki.log 2>&1 