# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/karma/run
#!/usr/bin/with-contenv bashio

echo "Starting Karma..."
exec karma --config.file=/etc/karma/karma.yml 