# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/karma/run
#!/usr/bin/with-contenv bashio

echo "Starting Karma..."
sleep 5
exec /usr/local/bin/karma --config.file=/etc/karma/karma.yml 