# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/karma/run
#!/command/with-contenv bashio

# Wait until Alertmanager API is reachable before starting Karma
for i in $(seq 1 30); do
    if curl -s "http://localhost:9093/api/v2/status" > /dev/null 2>&1; then
        echo "Alertmanager API is reachable (attempt $i)"
        break
    fi
    echo "Waiting for Alertmanager API... (attempt $i/30)"
    sleep 0.5
    if [ $i -eq 30 ]; then
        echo "Alertmanager API not reachable after 15s, proceeding anyway"
    fi
done

# Start karma
echo "Starting Karma..."
exec s6-setuidgid root /usr/local/bin/karma \
    --config.file=/etc/karma/karma.yml \
    >> /var/log/karma/karma.log 2>&1

# Check Prometheus health
if [ "$service_name" = "Prometheus" ]; then
    for n in {1..30}; do
        if curl -s http://localhost:9090/api/v1/targets |
           grep -q '"health":"up"'; then
            print_success "✅ Can scrape targets"
            return 0
        fi
        sleep 1
    done
    print_error "❌ Cannot scrape targets"
    return 1
fi 