# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/karma/run
#!/command/with-contenv bashio

# Start karma
exec s6-setuidgid root /usr/local/bin/karma \
    --config.file=/etc/karma/karma.yml 