# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/karma/run
#!/usr/bin/with-contenv bashio

echo "Starting Karma..."
exec s6-setuidgid root /usr/local/bin/karma \
    --config.file=/etc/karma/karma.yml \
    --alertmanager.uri=http://localhost:9093 