# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/karma/up
#!/command/with-contenv bashio

# Wait for Karma to be ready
for i in $(seq 1 30); do
    if curl -s "http://localhost:8080/health" > /dev/null; then
        echo "Karma is ready"
        exit 0
    fi
    echo "Waiting for Karma... (attempt $i/30)"
    sleep 0.5
done

echo "Karma failed to become ready after 15 seconds"
exit 1 