# prometheus-stack/rootfs/etc/s6-overlay/s6-rc.d/karma/up
#!/usr/bin/with-contenv bashio

echo "Karma service is up" 